# crontab
*/5 * * * * /bin/sh /var/www/storage/app/script/gera-logs.sh -e
